# .config/systemd/user/memory-bridge.service
[Unit]
Description=ANAQ Memory Bridge (port 9600)
After=network.target embedding-service.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
ExecStart=/home/fareez541/miniforge3/envs/agent0/bin/python /home/fareez541/.anaq/bridge/memory_bridge.py
WorkingDirectory=/home/fareez541/.anaq/bridge
Environment=HOME=/home/fareez541
Environment=PATH=/home/fareez541/miniforge3/envs/agent0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=ROCM_HOME=/opt/rocm
Environment=HSA_OVERRIDE_GFX_VERSION=11.0.0
Restart=always
RestartSec=3
MemoryMax=4G
OOMScoreAdjust=-500
StandardOutput=journal
StandardError=journal
SyslogIdentifier=memory-bridge

[Install]
WantedBy=default.target
